# playbooks/traefik/files/traefik.service.j2
[Unit]
Description=Traefik v3 Edge Router
After=network-online.target
Wants=network-online.target

[Service]
Type=notify
User=root
Group=root

# Service configuration
ExecStart=/usr/local/bin/traefik --configfile={{ traefik_config_dir }}/traefik.yml
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=traefik

[Install]
WantedBy=multi-user.target
